# supervisord.conf
[supervisord]
nodaemon = true

[program:postfix]
command = postfix start